# roles/npm-workstations/templates/etc/profile.d/node.sh.j2
# {{ ansible_banner }}

export PATH=${PATH}:{{ node_home }}/bin
